# .config/autostart/open_browser.desktop
[Desktop Entry]
Type=Application
Name=OpenBrowser
Exec=/home/pi/tw-stream/browser-start.sh &
X-GNOME-Autostart-enabled=true